# support/systemd/outbreak-running.target
[Unit]
Description=Running state, in which the Outbreak web app is running
Requires=outbreak.service

