# support/systemd/restoring.target
[Unit]
Description=Startup state, in which all that is running is webapp
Requires=restore-app.service

