# etc/cron.d/shapshot.cron
# Take a db snapshot every 15 minutes
*/15 * * * * bash /support/snapshot.sh
