# etc/profile.d/pythonpath.sh
#!/bin/bash
# Probably not necessary, really.
export PYTHONPATH="/usr/share"

